# node1-ansible/roles/operations_deployment/templates/clean-logs.service.j2
[Unit]
Description=Clean old journal logs

[Service]
Type=oneshot
ExecStart=/usr/bin/journalctl --vacuum-time=14d
